# spbot.service
[Unit]
Description=ScreenPlan QQ Bot - AI Agent
After=network.target screenplan.service

[Service]
Type=simple
User=screenplan
Group=screenplan
WorkingDirectory=/opt/spbot

Environment="SPBOT_NAPCAT_WS=ws://172.17.0.1:3001"
Environment="SPBOT_LISTEN_HOST=0.0.0.0"
Environment="SPBOT_LISTEN_PORT=3001"
Environment="SPBOT_ADMIN_TOKEN=8L4ADdOL4rx79GNhVnWfQ07sl-rhCZ-BgRj5lmOqbSc"
Environment="SPBOT_SCREENPLAN_HOST=http://localhost:5051"
Environment="SPBOT_LLM_BASE=https://api.deepseek.com/v1"
Environment="SPBOT_LLM_MODEL=deepseek-chat"
Environment="SPBOT_LLM_KEY=YOUR_DEEPSEEK_API_KEY"

ExecStart=/opt/screenplan-backend/venv/bin/python3 /opt/spbot/spbot.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=spbot

PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/spbot

[Install]
WantedBy=multi-user.target
